# deploy/pokestop-bot.service
[Unit]
Description=Yujin's Pokestop - Discord bot (bot.py, react-to-track auctions)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/low-sale-finder
ExecStart=/home/ubuntu/low-sale-finder/venv/bin/python bot.py
Restart=always
RestartSec=30
StandardOutput=append:/home/ubuntu/low-sale-finder/logs/bot.log
StandardError=append:/home/ubuntu/low-sale-finder/logs/bot.log

[Install]
WantedBy=multi-user.target
